# launch-agents/env.dotnet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>env.dotnet</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/launchctl</string>
		<string>setenv</string>
		<string>DOTNET_ROOT</string>
		<string>/opt/homebrew/opt/dotnet/libexec</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
